# family-dashboard.service

[Unit]
Description=Family Dashboard
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/family-dashboard
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=5000

[Install]
WantedBy=multi-user.target
